# daemon/systemd/asteroid-sensorlogd.service
[Unit]
Description=Sensor logger for AsteroidOS

[Service]
ExecStart=/usr/bin/sensorlogd

[Install]
WantedBy=basic.target
